# bin/selfhosting-ceo.service
[Unit]
Description=selfhosting.sh - CEO Agent
After=network.target

[Service]
Type=simple
ExecStart=/opt/selfhosting-sh/bin/run-agent.sh /opt/selfhosting-sh 3600
Restart=always
RestartSec=30
User=selfhosting
WorkingDirectory=/opt/selfhosting-sh
Environment=HOME=/root
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
